# tools/shadcn-comparison-gallery/ops/swiftcn-gallery.service
[Unit]
Description=Swiftcn shadcn visual comparison gallery
After=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/simon/github/swiftcn-ui/tools/shadcn-comparison-gallery
ExecStart=/home/simon/github/swiftcn-ui/tools/shadcn-comparison-gallery/ops/serve-gallery.sh
Restart=on-failure
RestartSec=2

[Install]
WantedBy=default.target
